# deploy/com.yourart.ai-governor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.yourart.ai-governor</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/ai-governor/deploy/start.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>StandardOutPath</key>
    <string>/var/log/ai-governor/launchd.log</string>

    <key>StandardErrorPath</key>
    <string>/var/log/ai-governor/launchd.err</string>

    <key>WorkingDirectory</key>
    <string>/opt/ai-governor</string>
</dict>
</plist>
